# tasks.crontab
0 0 * * * /usr/bin/php /var/www/html/cleanup.php >> /var/log/chat-site-cron.log 2>&1
